# web_server_backup.cron
00 7 * * * /usr/local/web-server-backup/backup.sh > /usr/local/web-server-backup/backup.log 2>&1
